# attendance-api.service
[Unit]
Description=Jinhe Attendance API Server
After=network.target

[Service]
Type=simple
User=www-data
ExecStart=/usr/bin/python3 /opt/attendance-api/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
